# monitorrent/rootfs/etc/services.d/nginx/run
#!/usr/bin/with-contenv bashio

# Wait for monitorrent to become available
bashio::net.wait_for 6687 localhost 120

bashio::log.info "Starting NGinx..."

# Disable HA Authentication if front door is open
if bashio::config.true 'leave_front_door_open'; then
    export DISABLE_HA_AUTHENTICATION=true
fi

exec nginx
